# rootfs/etc/services.d/access_control_server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Access Control
# Runs Access Control
# ==============================================================================
bashio::log.info "Starting Access Control server."
cd /opt/server/app || bashio::exit "Could not change dir."
exec python3 main.py
